# gateway.service
[Unit]
Description=Anthropic-OpenAI Gateway Service
Documentation=https://github.com/your-repo/anthropic-openai-gateway
After=network.target

[Service]
Type=simple
User=gateway
Group=gateway
WorkingDirectory=/opt/gateway

# 核心配置
ExecStart=/opt/gateway/gateway

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=600
StartLimitBurst=5

# 日志输出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gateway

# 资源限制
LimitNOFILE=65535
LimitNPROC=32768

# 环境变量（按需修改）
EnvironmentFile=-/etc/gateway/gateway.env
Environment="LISTEN_HOST=0.0.0.0"
Environment="LISTEN_PORT=3456"

# 推荐把 OpenAI 与模型配置放到 /etc/gateway/gateway.env：
# OPENAI_API_KEY=sk-xxx
# OPENAI_BASE_URL=https://api.openai.com/v1
# OPENAI_TIMEOUT_MS=120000
# MODELS_NEED_TRANSFORMATION=gpt-4.1,gpt-4o,gpt-4o-mini,gpt-4.1-mini,gpt-4.1-nano,gpt-5,o3,o3-mini,o4-mini

[Install]
WantedBy=multi-user.target
